# .crontab
Shell=/bin/zsh
0 7 * * *  gem update -V -y
0 8 * * *  brew update && brew upgrade
